# deployment_backups/20250829_190929/hooks_backup/shadowgit/shadowgit-systemd-service.service
# /etc/systemd/system/shadowgit-neural.service
# Shadowgit Neural-Enhanced File Monitoring Service
# 
# Installation:
#   sudo cp shadowgit-neural.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable shadowgit-neural.service
#   sudo systemctl start shadowgit-neural.service

[Unit]
Description=Shadowgit Neural File Monitoring with NPU/GNA Acceleration
Documentation=https://github.com/shadowgit/neural-docs
After=network.target multi-user.target

# Dependencies
Wants=network-online.target
After=network-online.target

# Ordering
Before=docker.service
After=basic.target

[Service]
Type=simple
User=shadowgit
Group=shadowgit
WorkingDirectory=/opt/shadowgit

# Environment Variables
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/opt/shadowgit:/opt/shadowgit/lib"

# OpenVINO Environment
Environment="OPENVINO_ROOT=/opt/intel/openvino"
Environment="LD_LIBRARY_PATH=/opt/intel/openvino/runtime/lib/intel64:/opt/intel/openvino/runtime/3rdparty/tbb/lib:$LD_LIBRARY_PATH"
Environment="PYTHONPATH=/opt/intel/openvino/python/python3.10:/opt/intel/openvino/tools:$PYTHONPATH"
Environment="CMAKE_PREFIX_PATH=/opt/intel/openvino/runtime/cmake"
Environment="ngraph_DIR=/opt/intel/openvino/runtime/cmake"
Environment="InferenceEngine_DIR=/opt/intel/openvino/runtime/cmake"
Environment="IE_PLUGINS_PATH=/opt/intel/openvino/runtime/lib/intel64"

# NPU Configuration
Environment="ZE_INTEL_NPU_LOGLEVEL=WARNING"
Environment="OV_NPU_LOG_LEVEL=WARNING"
Environment="NPU_COMPILER_TYPE=DRIVER"
Environment="NPU_PLATFORM=3720"
Environment="NPU_PROFILING=0"
Environment="NPU_PRINT_PROFILING=0"
Environment="NPU_DPU_PROFILING=0"
Environment="NPU_DMA_PROFILING=0"
Environment="NPU_SW_PROFILING=0"
Environment="NPU_BYPASS_UMD_CACHING=0"

# GNA Configuration  
Environment="GNA_DEVICE_MODE=GNA_HW"
Environment="GNA_PRECISION=I8"
Environment="GNA_PWL_UNIFORM_DESIGN=YES"
Environment="GNA_FIRMWARE_MODEL_IMAGE=STREAMING"
Environment="GNA_EXEC_TARGET=GNA_TARGET_3_5"
Environment="GNA_COMPILE_TARGET=GNA_TARGET_3_5"
Environment="GNA_LIB_N_THREADS=1"

# Performance Tuning
Environment="OMP_NUM_THREADS=4"
Environment="MKL_NUM_THREADS=4"
Environment="OV_CPU_DENORMALS_OPTIMIZATION=YES"
Environment="OV_CPU_SPARSE_WEIGHTS_DECOMPRESSION_RATE=0.5"
Environment="TBB_MALLOC_USE_HUGE_PAGES=1"

# Python Optimization
Environment="PYTHONOPTIMIZE=1"
Environment="PYTHONDONTWRITEBYTECODE=1"

# Shadowgit Configuration
Environment="SHADOWGIT_MODE=production"
Environment="SHADOWGIT_NEURAL=enabled"
Environment="SHADOWGIT_POWER_MODE=balanced"
Environment="SHADOWGIT_BATCH_SIZE=32"
Environment="SHADOWGIT_BATCH_WINDOW_MS=500"
Environment="SHADOWGIT_GNA_ALWAYS_ON=true"
Environment="SHADOWGIT_NPU_PRIORITY=high"

# Security
Environment="SHADOWGIT_SECURE_MODE=true"
Environment="SHADOWGIT_SANDBOX=enabled"

# Main execution command
ExecStartPre=/opt/shadowgit/scripts/check_hardware.sh
ExecStart=/usr/bin/python3 /opt/shadowgit/shadowgit_watcher_enhanced.py \
    --watch /home/shadowgit/projects \
    --power-mode balanced \
    --batch-size 32 \
    --batch-window 500

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
LimitCORE=infinity
LimitMEMLOCK=infinity

# Memory management
MemoryMax=4G
MemorySwapMax=1G
MemoryLow=512M

# CPU management
CPUWeight=100
CPUQuota=400%
AllowedCPUs=0-11

# IO management
IOWeight=100
IOReadBandwidthMax=/dev/nvme0n1 200M
IOWriteBandwidthMax=/dev/nvme0n1 100M

# Security hardening
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/shadowgit /var/log/shadowgit /tmp
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
NoNewPrivileges=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Device access for NPU/GNA
DeviceAllow=/dev/dri/renderD128 rw
DeviceAllow=/dev/dri/card0 rw
DeviceAllow=/dev/gna0 rw
DeviceAllow=/dev/accel/accel0 rw
PrivateDevices=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shadowgit-neural

# Performance monitoring
CollectMode=inactive-or-failed

# Watchdog
WatchdogSec=300
NotifyAccess=main

# Process tracking
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
SendSIGKILL=yes

# Nice level (priority)
Nice=-10

# Task scheduling
TasksMax=512

[Install]
WantedBy=multi-user.target
Alias=shadowgit.service

# Additional dependencies for installation
Also=shadowgit-gna-monitor.timer