# agents/services/claude-agents.service
[Unit]
Description=Claude Agent Communication System
After=network.target

[Service]
Type=forking
User=ubuntu
WorkingDirectory=${CLAUDE_PROJECT_ROOT:-$(dirname "$0")/../../}agents
ExecStart=${CLAUDE_AGENTS_ROOT:-$(dirname "$0")}/../BRING_ONLINE.sh
ExecStop=/usr/bin/pkill -f unified_agent_runtime
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
